# deploy/systemd/tatuscan.socket
[Unit]
Description=TatuScan Server Socket
PartOf=tatuscan.service

[Socket]
ListenStream=8040
BindIPv6Only=both

[Install]
WantedBy=sockets.target